# init/gtfs-editor.conf
# Upstart script to start GTFS Editor

env user=matthewc
# configure port and bind address in application.cong

# application home directory
chdir /home/matthewc/conveyal/gtfs-editor

start on runlevel [3]

pre-start script
  # We set up an exclusive temporary directory so that Geotools will play nice with other Geotools apps on the system,
  # which may be running as different users
  # see http://www.indicatrix.org/2014/10/20/using-geotools-with-multiple-user-accounts/

  # The environment stuff with writing to a file is a bit annoying; see http://stackoverflow.com/questions/12610335
  TMPDIR=`mktemp -d`
  chown $user $TMPDIR
  echo "TMPDIR=$TMPDIR" > "/tmp/$UPSTART_JOB"
end script

script
  . "/tmp/$UPSTART_JOB"
  su -m -c "exec play run --%production" $user
end script

post-stop script
  # Clear out our temporary directory
  . "/tmp/$UPSTART_JOB"
  rm -rf "$TMPDIR"
  rm "/tmp/$UPSTART_JOB"
end script
